# deployment/gunicorn.service
[Unit]
Description=Gunicorn instance to serve shopping-react backend (FastAPI)
After=network.target

[Service]
User=ai-tracks-shopping-react
Group=ai-tracks-shopping-react
WorkingDirectory=/home/ai-tracks-shopping-react/htdocs/shopping-react.ai-tracks.com/backend
Environment="PATH=/home/ai-tracks-shopping-react/htdocs/shopping-react.ai-tracks.com/backend/.venv/bin:/usr/local/bin:/usr/bin:/bin"

# 使用 gunicorn + uvicorn workers（推荐方案，原生支持 ASGI）
# 安装 gunicorn: uv add gunicorn
ExecStart=/home/ai-tracks-shopping-react/htdocs/shopping-react.ai-tracks.com/backend/.venv/bin/gunicorn \
    app.main:app \
    -w 8 \
    -k uvicorn.workers.UvicornWorker \
    -b 127.0.0.1:8096 \
    --access-logfile /var/log/gunicorn/shopping-react-access.log \
    --error-logfile /var/log/gunicorn/shopping-react-error.log \
    --log-level info

Restart=always
KillSignal=SIGTERM
Type=notify
NotifyAccess=all

[Install]
WantedBy=multi-user.target

